# launchagents/com.taobinxian.homework-tutor-8787.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.taobinxian.homework-tutor-8787</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>proxy.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/homework-tutor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>8787</string>
    <key>BIND</key>
    <string>0.0.0.0</string>
    <key>STATIC_DIR</key>
    <string>/ABSOLUTE/PATH/TO/homework-tutor</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/homework-tutor-8787.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/homework-tutor-8787.log</string>
</dict>
</plist>
